# root/etc/services.d/servicename/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc /bin/bash -c \
    "/app/sbnc/bin/sbnc --config /config --foreground"
